# scripts/crontab
# 社交媒体自动发布定时任务
# 每天上午 9:00 执行社交媒体发布任务
0 9 * * * cd /path/to/awesome-crypto-tutorial.github.io && /usr/bin/ruby scripts/social-media-bot.rb >> logs/social-media-bot.log 2>&1
