# scripts/crontab
# 社交媒体自动发布定时任务
# 每天上午 9:00 执行社交媒体发布任务
0 9 * * * /bin/bash -l -c "cd /root/awesome-crypto-tutorial.github.io && . scripts/run-bot.sh >> logs/social-media-bot.log 2>&1"
